# quassel-core/create-cert.conf
[program:create-cert]
command = /home/create-cert.sh
autostart = true
autorestart = false
startsecs = 0
user = nobody
group = users
#priority = 999
