# quassel-core/quassel-core.conf
[program:quasselcore]
command = /usr/sbin/quasselcore --listen=0.0.0.0 --configdir=/config/
#command = /home/start-quassel.sh
autostart = false
autorestart = true
user = nobody
group = users
