# firefox-sync/ffsync.conf
[program:ffsync]
command=/opt/mozilla-firefox-sync-server/local/bin/gunicorn --paste /config/syncserver.ini
directory=/home
autorestart=true
user=nobody
group=users
#environment=STNORESTART="1"
