# firefox-sync/fxa-local-dev.conf
[program:fxa-local-dev]
command=/home/launch_fxa-local-dev.sh
directory=/home
autorestart=true
user=nobody
group=users
environment=NVM_DIR="/home/.nvm"
