# syncthing/prep-service.conf
[program:prep-service]
command = /home/prep-service.sh
autostart = true
autorestart = false
startsecs = 0
user = nobody
group = users
#priority = 999
